# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── Xvfb — shared virtual display ─────────────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX
environment=DISPLAY=""
autostart=true
autorestart=true
priority=10
startsecs=1
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── x11vnc — VNC server on :99, internal port 5900 ───────────
[program:x11vnc]
command=x11vnc -display :99 -rfbport 5900 -nopw -forever -shared -noxdamage -quiet -wait_ui 5000
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── noVNC — browser-accessible VNC on port 6080 ──────────────
[program:novnc]
command=websockify --web /usr/share/novnc/ 6080 localhost:5900
autostart=true
autorestart=true
priority=30
startsecs=4
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Fantoma HTTP server ────────────────────────────────────────
[program:fantoma]
command=python3 /app/server.py
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/fantoma.log
stderr_logfile=/var/log/supervisor/fantoma.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
